# deployment/fixer-ai.service
[Unit]
Description=Fixer AI Repair Agent
After=network.target

[Service]
Type=simple
User=fixer
WorkingDirectory=/opt/fixer-ai
ExecStart=/usr/bin/python3 /opt/fixer-ai/main.py --sms-daemon
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=fixer-ai
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
